# dictation.service
[Unit]
Description=Dictation hotkey daemon

[Service]
ExecStart=/usr/bin/python3 /home/vkommera/Documents/Projects/dictation/main.py
Restart=on-failure
Environment=WAYLAND_DISPLAY=wayland-0
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=default.target
